# simple-usb-notifier/etc/udev/rules.d/99-usb-notification.rules
ACTION=="add|remove", SUBSYSTEM=="usb", ENV{DISPLAY}=":0", RUN+="/bin/bash -c 'susb-notifier $env{DEVNAME} $env{ACTION}'"
